# broadlinkac/files/etc/uci-defaults/99-broadlinkac
#!/bin/sh
# First-run: generate config.json from UCI settings
[ -f /root/.ac_controller/config.json ] && exit 0

mkdir -p /root/.ac_controller/logs

uci -q get broadlinkac.@settings[0] > /dev/null 2>&1 || exit 0

get() { uci -q get broadlinkac.@settings[0].$1 || echo ""; }
getb() { val=$(uci -q get broadlinkac.@settings[0].$1); [ "$val" = "1" ] && echo "true" || echo "false"; }

cat > /root/.ac_controller/config.json << 'JSONEOF'
{
  "api_key": "APIKEY",
  "qw_host": "QWHOST",
  "baidu_key": "BDKEY",
  "weather_provider": "WP",
  "weather_provider_set": true,
  "location": {
    "lat": LAT,
    "lon": LON,
    "name": "LOCNAME"
  },
  "typhoon_provider": "TP",
  "typhoon_ac_off": ACOFF,
  "appearance_mode": "system",
  "devices": {},
  "current_device_mac": ""
}
JSONEOF

sed -i "s/APIKEY/$(get api_key)/" /root/.ac_controller/config.json
sed -i "s|QWHOST|$(get qw_host)|" /root/.ac_controller/config.json
sed -i "s/BDKEY/$(get baidu_key)/" /root/.ac_controller/config.json
sed -i "s/WP/$(get weather_provider)/" /root/.ac_controller/config.json
sed -i "s/LAT/$(get location_lat)/" /root/.ac_controller/config.json
sed -i "s/LON/$(get location_lon)/" /root/.ac_controller/config.json
sed -i "s/LOCNAME/$(get location_name)/" /root/.ac_controller/config.json
sed -i "s/TP/$(get typhoon_provider)/" /root/.ac_controller/config.json
sed -i "s/ACOFF/$(getb typhoon_ac_off)/" /root/.ac_controller/config.json

echo "BroadlinkAC config initialized"
